# istoreos-files/files/etc/uci-defaults/99_fullcone_nat_fw4
#!/bin/sh

[ "0" = "$(uci -q get 'firewall.@defaults[0].fullcone')" ] && exit 0

uci -q batch <<-EOF
	set firewall.@defaults[0].fullcone='2'
	commit firewall
EOF

exit 0
